# templates/monitoring.service.j2
[Unit]
Description=Monitoring script for the test process
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/monitor.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
